# deploy/db-backup.service
[Unit]
Description=Database backup service

[Service]
Type=oneshot
User=app
WorkingDirectory=/opt/app
ExecStart=/opt/app/db-backup.sh
EnvironmentFile=/opt/app/app.env
StandardOutput=journal
StandardError=journal
